# config/includes.chroot/etc/systemd/system/tuxedo.service
[Unit]
Description=Enable TUXEDO power savings

[Service]
Type=oneshot
ExecStart=/bin/tuxedopowersave

[Install]
WantedBy=multi-user.target
